# dockerfiles/crontab
PYTHONPATH=/srv/cache_school_data

* * * * * root echo hello $PYTHONPATH > /proc/1/fd/1 2>/proc/1/fd/2
