# dockerfiles/crontab
PYTHONPATH=/srv/cache_school_data

* * * * * root python3 /srv/cache_school_data/cache_school_data/main.py > /proc/1/fd/1 2>/proc/1/fd/2

